# supervisor/dropbear.conf
[program:dropbear]
command=/usr/sbin/dropbear -F
